# linux/udev/rules.d/80-yubikey-removed.rules
ACTION=="remove", SUBSYSTEMS=="usb", ATTRS{idVendor}=="1050", RUN+="/home/icetan/.local/bin/lock.sh"
