# linux/udev/rules.d/90-keyboard-layout.rules
ACTION=="add", SUBSYSTEMS=="input", RUN+="/home/icetan/.local/bin/upkeys", OWNER="icetan"
